# launchd/com.user.finnews.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.finnews</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/och/.venvs/fin-news-collector/bin/python3</string>
    <string>/Users/och/코딩/harness practice/fin-news-collector/run.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/och/코딩/harness practice/fin-news-collector</string>
  <key>StandardOutPath</key>
  <string>/tmp/finnews.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/finnews.err</string>
</dict>
</plist>
